# deployment/systemd/tailtown-frontend.service
[Unit]
Description=Tailtown Frontend (Nginx serves static files)
After=network.target nginx.service
Wants=nginx.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecReload=/usr/bin/systemctl reload nginx

[Install]
WantedBy=multi-user.target
